# infra/voice-agent.service
[Unit]
Description=LevelOne Voice Agent Worker
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/voice-agent
EnvironmentFile=/opt/voice-agent/.env
ExecStart=/opt/voice-agent/.venv/bin/voice-worker start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
